# cookbooks/solr/templates/default/solr.erb
#!/bin/sh
#
# solr       Starts the solr service
#
# chkconfig:   345 05 95
# description:  
# ....

cd /usr/local/solr/apache-solr-3.6.0/example
java -jar start.jar &